# bin/99-camera.rules
#to listen to usb: udevadm monitor -p -u -s usb

#run when usb added the script should determine whether or not /dev/video0 exists
ACTION=="add", SUBSYSTEMS=="usb", RUN+="/home/pi/videoexists.sh"


